# etc/dist/zoia.dist.service
[Unit]
Description=ZOIA Web Server ({site_id})
Documentation=https://zoiajs.org
After=network.target

[Service]
Environment=
Type=simple
User={user}
Group={group}
ExecStart=/usr/bin/node {root}/build/bin/server.js
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=zoia_{site_id}

[Install]
WantedBy=multi-user.target